# turtlebot_bringup/upstart/ps3joy.conf
start on runlevel [2345]

respawn

exec /opt/ros/electric/stacks/joystick_drivers/ps3joy/ps3joy.py
