# turtlebot_bringup/upstart/turtlebot.conf
description	"bringup turtlebot"

start on net-device-up IFACE=wlan0
stop on net-device-down IFACE=wlan0

console output
respawn

exec turtlebot-start
pre-stop exec turtlebot-stop
