# turtlebot_bringup/upstart/turtlebot.conf
description	"bringup turtlebot"

start on turtlebot-session-begin
stop on turtlebot-session-end

console output
respawn

exec turtlebot-start $IFACE
pre-stop exec turtlebot-stop
